# daemon/src/mydaemon.service
[Unit]
Description=MyDaemon

[Service]
Type=simple
ExecStart=/usr/bin/mydaemond

[Install]
WantedBy=multi-user.target

