# daemon/src/mydaemon.service
[Unit]
Description=MyDaemon

[Service]
Type=simple
ExecStart=/usr/bin/mydaemon

[Install]

